# package/independent/base-settings-xwrt-xf-949/files/rc.local

env -i ACTION="ifchange" INTERFACE="br-lan" LINK=1 /sbin/hotplug-call iface

echo "1" >/sys/kernel/debug/ieee80211/phy0/mt76/edcca

exit 0
